# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── nginx ──────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx_err.log

# ── FastAPI backend ────────────────────────────────────────
[program:backend]
command=uvicorn main:app --host 127.0.0.1 --port 8000 --workers 2
directory=/app/backend
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/backend.log
stderr_logfile=/var/log/supervisor/backend_err.log
environment=
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s",
    ADMIN_TOKEN="%(ENV_ADMIN_TOKEN)s",
    ENVIRONMENT="production"

# ── Telegram Bot (polling) ────────────────────────────────
[program:bot]
command=python main.py
directory=/app/bot
autostart=true
autorestart=true
priority=30
startretries=5
stdout_logfile=/var/log/supervisor/bot.log
stderr_logfile=/var/log/supervisor/bot_err.log
environment=
    BOT_TOKEN="%(ENV_BOT_TOKEN)s",
    BACKEND_URL="http://127.0.0.1:8000",
    REDIS_URL="%(ENV_REDIS_URL)s",
    ENVIRONMENT="production"

# ── Admin Panel (Next.js standalone) ─────────────────────
[program:admin]
command=node server.js
directory=/app/admin
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/admin.log
stderr_logfile=/var/log/supervisor/admin_err.log
environment=
    PORT="3000",
    HOSTNAME="127.0.0.1",
    NODE_ENV="production"

# ── Mini App (Next.js standalone) ────────────────────────
[program:mini]
command=node server.js
directory=/app/mini
autostart=true
autorestart=true
priority=40
stdout_logfile=/var/log/supervisor/mini.log
stderr_logfile=/var/log/supervisor/mini_err.log
environment=
    PORT="3001",
    HOSTNAME="127.0.0.1",
    NODE_ENV="production"
